# profile.d/default_editor
#!/bin/sh

# Define the preferred order of editors
for VISUAL in vim nano vi; do
    which "$VISUAL" >/dev/null 2>&1 && break
done
export EDITOR="$VISUAL" # pass(1) relies on EDITOR
