# profile.d/undistract-me
[ -f  /etc/profile.d/undistract-me.sh ] \
    && . /etc/profile.d/undistract-me.sh
